# rootfs/etc/services.d/openvpn/run
#!/usr/bin/with-contenv sh
/usr/local/bin/ovpn_run